# usr/lib/systemd/system/biglinux-key.service
[Unit]
Description=copy keys to /etc/pacman.d/gnupg/
After=pacman-init.service


[Service]
Type=simple
ExecStart=/usr/bin/biglinux-key.sh


[Install]
WantedBy=multi-user.target 
